# peat-node-pkg/debian/peat-node.service
[Unit]
Description=Defense Unicorns peat-node (CRDT mesh sync daemon)
Documentation=https://github.com/defenseunicorns/peat-node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=peat
Group=peat

# Load all PEAT_NODE_* variables from the conffile
EnvironmentFile=/etc/peat-node/peat-node.env

# Split PEAT_NODE_PEERS into --peer flags if set.
# peat-node accepts repeated --peer flags; we pass them via the env var
# natively if the binary supports PEAT_NODE_PEERS, otherwise use ExecStartPre
# to build a wrapper. Check peat-node --help for native support.
ExecStart=/usr/sbin/peat-node

WorkingDirectory=/var/lib/peat-node

# Restart policy: always restart on crash; not on clean SIGTERM (package remove)
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
# peat-node needs to write its data dir
ReadWritePaths=/var/lib/peat-node
# Iroh QUIC needs UDP; gRPC needs TCP
RestrictAddressFamilies=AF_INET AF_INET6
AmbientCapabilities=CAP_NET_RAW
# Prevent privilege escalation via setuid binaries
LockPersonality=true
RestrictRealtime=true

# Logging — follow with: journalctl -u peat-node -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=peat-node

[Install]
WantedBy=multi-user.target
