# peat-node-pkg/debian/peat-node.service
[Unit]
Description=Defense Unicorns peat-node (CRDT mesh sync daemon)
Documentation=https://github.com/defenseunicorns/peat-node
After=network-online.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service

[Service]
Type=simple
User=peat
Group=peat

# Load all PEAT_NODE_* variables from the conffile
EnvironmentFile=/etc/peat-node/peat-node.env

ExecStart=/usr/sbin/peat-node

WorkingDirectory=/var/lib/peat-node

# Restart policy: always restart on crash; not on clean SIGTERM (package remove)
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL

# Security hardening
NoNewPrivileges=true
# PrivateTmp would isolate /tmp and break mDNS socket access — must be off
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
# peat-node needs to write its data dir
ReadWritePaths=/var/lib/peat-node
# AF_NETLINK is required for multicast group membership (mDNS)
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK
# CAP_NET_RAW and CAP_NET_BROADCAST are needed to join multicast groups
# for Iroh's local mDNS peer discovery
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
# Prevent privilege escalation via setuid binaries
LockPersonality=true
RestrictRealtime=true

# Logging — follow with: journalctl -u peat-node -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=peat-node

[Install]
WantedBy=multi-user.target
